# dozzle/rootfs/etc/services.d/dozzle/finish
#!/bin/shio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Dozzle
# Service shutdown script
# ==============================================================================

bashio::log.info "Stopping Dozzle service..."

# Exit code for S6
exit 0 