# dozzle/rootfs/etc/services.d/dozzle/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Dozzle
# ==============================================================================

# ==============================================================================
# FINISH LOGIC
# ------------------------------------------------------------------------------
main() {
    bashio::log.info "Stopping Dozzle..."
}

# ==============================================================================
# EXECUTE LOGIC
# ------------------------------------------------------------------------------
main "$@" 