# dozzle/rootfs/etc/services.d/dozzle/finish
#!/bin/bash
# ==============================================================================
# Home Assistant Add-on: Dozzle
# ==============================================================================

# Simple logging function
log_message() {
    echo "$(date '+%Y-%m-%d %H:%M:%S') $1"
}

# Get exit code from service
EXIT_CODE=$1

log_message "Dozzle service exited with code ${EXIT_CODE}"

# Forward exit code
exit ${EXIT_CODE} 