# dozzle/rootfs/etc/services.d/dozzle/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Dozzle
# Start Dozzle service
# ==============================================================================

bashio::log.info "Starting Dozzle service via S6..."

# Execute main script with proper environment
cd /
exec /run.sh 