# dozzle/rootfs/etc/services.d/dozzle/run
#!/bin/bash
# ==============================================================================
# Home Assistant Add-on: Dozzle
# ==============================================================================

# Execute the main run.sh script
exec /run.sh 