# dozzle/rootfs/etc/services.d/dozzle/run
#!/bin/bash
# ==============================================================================
# Home Assistant Add-on: Dozzle
# ==============================================================================

# Simple logging function
log_message() {
    echo "$(date '+%Y-%m-%d %H:%M:%S') $1"
}

log_message "Starting Dozzle..."

# Source configuration
CONFIG_PATH=/data/options.json

# Get configuration values
LOG_LEVEL=$(jq -r '.log_level // "info"' $CONFIG_PATH)
EXTERNAL_ACCESS=$(jq -r '.external_access // "false"' $CONFIG_PATH)
AGENT_MODE=$(jq -r '.agent_mode // "false"' $CONFIG_PATH)

# Set Dozzle options
if [ "$EXTERNAL_ACCESS" = "true" ]; then
    # Use 0.0.0.0 to allow external access
    log_message "External access enabled"
    DOZZLE_OPTS="--addr 0.0.0.0:8080 --base / --level ${LOG_LEVEL}"
else
    # Use 127.0.0.1 for internal access only (ingress)
    log_message "External access disabled, using ingress only"
    DOZZLE_OPTS="--addr 127.0.0.1:8080 --base / --level ${LOG_LEVEL}"
fi

# Add agent mode if enabled
if [ "$AGENT_MODE" = "true" ]; then
    DOZZLE_OPTS="${DOZZLE_OPTS} --agent"
    log_message "Agent mode enabled"
fi

log_message "Dozzle options: ${DOZZLE_OPTS}"

# Run Dozzle
exec /usr/local/bin/dozzle ${DOZZLE_OPTS} 