# dozzle/rootfs/etc/services.d/dozzle/run
#!/bin/sh
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Dozzle
# Start Dozzle service
# ==============================================================================

# Define paths
BASHIO_BIN="/usr/bin/bashio"

# Check if bashio is available
if [ -x "${BASHIO_BIN}" ]; then
    # Use bashio if available
    log_info() { ${BASHIO_BIN} log.info "$1"; }
    log_error() { ${BASHIO_BIN} log.error "$1"; }
else
    # Fallback to echo if bashio is not available
    echo "Warning: bashio not found at ${BASHIO_BIN}, using fallback logging"
    log_info() { echo "[INFO] $1"; }
    log_error() { echo "[ERROR] $1"; }
fi

log_info "Starting Dozzle service via S6..."

# Check if run.sh exists and is executable
if [ ! -x /run.sh ]; then
    log_error "run.sh script not found or not executable!"
    exit 1
fi

# Execute main script with proper environment
cd /
exec /run.sh 