# dozzle/rootfs/etc/services.d/nginx/finish
#!/bin/sh
# ==============================================================================
# Home Assistant Add-on: Dozzle
# Take down the S6 supervision tree when Nginx fails
# ==============================================================================

# Define colors
RED='\033[0;31m'
GREEN='\033[0;32m'
YELLOW='\033[0;33m'
BLUE='\033[0;34m'
NC='\033[0m' # No Color

# Default log level (can be overridden by environment variable)
LOG_LEVEL="${LOG_LEVEL:-info}"

# Define log level values
LOG_LEVEL_DEBUG=3
LOG_LEVEL_INFO=2
LOG_LEVEL_WARNING=1
LOG_LEVEL_ERROR=0

# Set current log level based on configuration
CURRENT_LOG_LEVEL=$LOG_LEVEL_INFO
case "${LOG_LEVEL}" in
    "debug")
        CURRENT_LOG_LEVEL=$LOG_LEVEL_DEBUG
        ;;
    "info")
        CURRENT_LOG_LEVEL=$LOG_LEVEL_INFO
        ;;
    "error")
        CURRENT_LOG_LEVEL=$LOG_LEVEL_ERROR
        ;;
esac

# Use colored echo for logging with timestamp and level filtering
log_debug() { 
    if [ $CURRENT_LOG_LEVEL -ge $LOG_LEVEL_DEBUG ]; then
        timestamp=$(date +"%Y-%m-%d %H:%M:%S")
        printf "${BLUE}[DEBUG]${NC} ${BLUE}%s${NC} %s\n" "${timestamp}" "$1"
    fi
}

log_info() { 
    if [ $CURRENT_LOG_LEVEL -ge $LOG_LEVEL_INFO ]; then
        timestamp=$(date +"%Y-%m-%d %H:%M:%S")
        printf "${GREEN}[INFO]${NC} ${BLUE}%s${NC} %s\n" "${timestamp}" "$1"
    fi
}

log_warning() { 
    if [ $CURRENT_LOG_LEVEL -ge $LOG_LEVEL_WARNING ]; then
        timestamp=$(date +"%Y-%m-%d %H:%M:%S")
        printf "${YELLOW}[WARNING]${NC} ${BLUE}%s${NC} %s\n" "${timestamp}" "$1"
    fi
}

log_error() { 
    if [ $CURRENT_LOG_LEVEL -ge $LOG_LEVEL_ERROR ]; then
        timestamp=$(date +"%Y-%m-%d %H:%M:%S")
        printf "${RED}[ERROR]${NC} ${BLUE}%s${NC} %s\n" "${timestamp}" "$1"
    fi
}

# Get exit codes
EXIT_CODE_CONTAINER=0
if [ -f /run/s6-linux-init-container-results/exitcode ]; then
    EXIT_CODE_CONTAINER=$(cat /run/s6-linux-init-container-results/exitcode)
fi
EXIT_CODE_SERVICE="${1:-0}"
EXIT_CODE_SIGNAL="${2:-0}"

log_info "Service nginx exited with code ${EXIT_CODE_SERVICE} (by signal ${EXIT_CODE_SIGNAL})"

if [ "${EXIT_CODE_SERVICE}" -eq 0 ] || [ "${EXIT_CODE_CONTAINER}" -gt 0 ]; then
    log_info "Service restart not needed"
    exit 0
fi

log_info "Restarting service..."
s6-svc -r /var/run/s6/services/nginx 