# dozzle/rootfs/etc/services.d/nginx/run
#!/bin/sh
# ==============================================================================
# Home Assistant Add-on: Dozzle
# Runs the Nginx daemon for ingress support
# ==============================================================================

# Simple logging function
log_info() {
    timestamp=$(date +"%Y-%m-%d %H:%M:%S")
    echo "[INFO] ${timestamp} $1"
}

log_info "Starting NGinx for ingress support..."

# Run nginx
exec nginx -g "daemon off;" 