# dozzle/rootfs/etc/services.d/nginx/run
#!/bin/bash
# ==============================================================================
# Home Assistant Add-on: Dozzle
# Runs the Nginx daemon for ingress support
# ==============================================================================

# Simple logging function
log_message() {
    echo "$(date '+%Y-%m-%d %H:%M:%S') $1"
}

log_message "Starting Nginx..."

# Run nginx
exec nginx -g "daemon off;" 