# dozzle/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Dozzle
# Runs the Nginx daemon for ingress support
# ==============================================================================

# Use simple echo for logging
log_info() { echo "[INFO] $1"; }
log_warning() { echo "[WARNING] $1"; }

log_info "Starting NGinx for ingress support..."

# Run nginx
exec nginx -g "daemon off;" 